# scripts/launchd/com.copypro.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Auto-starts the CopyPro backend at login and restarts it if it crashes.
     Install: see scripts/launchd/README.md. Paths assume Neo's eli account. -->
<plist version="1.0"><dict>
  <key>Label</key><string>com.copypro.backend</string>
  <key>ProgramArguments</key><array>
    <string>/bin/bash</string>
    <string>/Users/eli/srv/harcourts-listings/run-copypro-backend.sh</string>
  </array>
  <key>WorkingDirectory</key><string>/Users/eli/srv/harcourts-listings</string>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>ThrottleInterval</key><integer>10</integer>
  <key>StandardOutPath</key><string>/Users/eli/srv/harcourts-listings/logs/daemon-backend.log</string>
  <key>StandardErrorPath</key><string>/Users/eli/srv/harcourts-listings/logs/daemon-backend.log</string>
</dict></plist>
